# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:ntpd]
command=/usr/sbin/ntpd -d
autorestart=true

[program:mysql]
command=/usr/bin/mysqld_safe
autorestart=true
priority=10

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
autorestart=true
priority=10

[program:keystone]
command=/usr/bin/keystone-all
autorestart=true
priority=20

[program:nova-api]
command=/usr/bin/nova-api
autorestart=true
priority=30

[program:nova-scheduler]
command=/usr/bin/nova-scheduler
autorestart=true
priority=30

[program:nova-cert]
command=/usr/bin/nova-cert
autorestart=true
priority=30

[program:nova-consoleauth]
command=/usr/bin/nova-consoleauth
autorestart=true
priority=30

[program:nova-conductor]
command=/usr/bin/nova-conductor
autorestart=true
priority=30

[program:glance-registry]
command=/usr/bin/glance-registry
autorestart=true
priority=30

[program:glance-api]
command=/usr/bin/glance-api
autorestart=true
priority=30

#[program:cinder-api]
#command=/usr/bin/cinder-api
#autorestart=true
#priority=30

#[program:cinder-scheduler]
#command=/usr/bin/cinder-scheduler
#autorestart=true
#priority=30

[program:memcached]
command=/usr/bin/memcached -u nova
autorestart=true
priority=10

[program:apache2]
command=/usr/sbin/apache2 -DFOREGROUND
autorestart=true
priority=40

#[program:docker]
#command=/usr/bin/wrapdocker
#autorestart=true





